# hub-gateway.service
[Unit]
Description=Home Services Hub Gateway
After=network.target

[Service]
Type=simple
User=myhomeservices
Group=myhomeservices
WorkingDirectory=/home/myhomeservices/hub
ExecStart=/home/myhomeservices/hub/hub-HomeServicesHub
ReadWritePaths=/home/myhomeservices/hub
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/myhomeservices/hub

[Install]
WantedBy=multi-user.target
